# charts/cron/crontab
* * * * * /scripts/trigger_xyz.sh >> /var/log/cron.log
